# root_ntp_nat.cron
# this is for interface em1
@reboot /sbin/iptables -w -t nat -A POSTROUTING -p udp --sport  123 -o em1 -j SNAT --to-source :1123     > /tmp/iptables.postroute 2>&1
@reboot /sbin/iptables -w -t nat -A PREROUTING  -p udp --dport 1123 -i em1 -j DNAT --to-destination :123 >  /tmp/iptables.preroute 2>&1
